# rootfs/etc/s6-overlay/s6-rc.d/03-initialise-gain/up
#!/bin/sh
exec /etc/s6-overlay/scripts/03-initialise-gain
