# rootfs/etc/s6-overlay/s6-rc.d/graphs_7d-1y/run
#!/bin/sh
exec /etc/s6-overlay/scripts/graphs_7d-1y
