# rootfs/etc/s6-overlay/s6-rc.d/graphs_1h-24h/run
#!/bin/sh
exec /etc/s6-overlay/scripts/graphs_1h-24h
