# rootfs/etc/s6-overlay/s6-rc.d/lighttpd/run
#!/bin/sh
exec /etc/s6-overlay/scripts/lighttpd
