# rootfs/etc/s6-overlay/s6-rc.d/telegraf_socat_vrs_json/run
#!/bin/sh
exec /etc/s6-overlay/scripts/telegraf_socat_vrs_json
