# rootfs/etc/s6-overlay/s6-rc.d/02-sanity-check/up
#!/bin/sh
exec /etc/s6-overlay/scripts/02-sanity-check
